# ops/systemd/supermarks-backend.service
[Unit]
Description=SuperMarks local production backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/repos/SuperMarks
ExecStart=%h/repos/SuperMarks/scripts/prod-backend.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
